# files/etc/services.d/inotifyd/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root /sbin/inotifyd /usr/bin/smokeping-restart.sh /config:c
